# sv/nginx/run
#!/bin/sh

chown -R www-data:www-data /var/www

exec 2>&1
exec /usr/sbin/nginx -g "daemon off;"
